# .deploy/websocket.conf
[program:websockets]
command=php /var/www/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=www-data
